# rc.local
sudo python /home/pi/Temperature/main.py;&